# deploy/aws/anavet-backend.service
# Servicio systemd del backend de ANA-vet.
#
# Importa sobre todo en el Learner Lab: cuando la sesión caduca la
# instancia se detiene, y al reanudarla systemd levanta el backend solo.
# Sin esto habría que entrar por SSH a arrancarlo a mano cada vez.
#
# Instalar en: /etc/systemd/system/anavet-backend.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now anavet-backend

[Unit]
Description=ANA-vet API (Express)
After=network.target mariadb.service mysqld.service
Wants=mariadb.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/anavet/backend/clinica-vet-backend
EnvironmentFile=/opt/anavet/backend/clinica-vet-backend/.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5

# La BD puede tardar en aceptar conexiones tras arrancar la instancia;
# que reintente sin rendirse en vez de quedarse en failed.
StartLimitIntervalSec=0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=anavet

[Install]
WantedBy=multi-user.target
